# docker/panel/supervisord.conf
[supervisord]
nodaemon = true
logfile  = /dev/stdout
logfile_maxbytes = 0
pidfile  = /tmp/supervisord.pid
user     = root

[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0700

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:ct-admin]
command          = bun run /opt/cool-tunnel/operator/src/index.ts admin serve
directory        = /opt/cool-tunnel/operator
autostart        = true
autorestart      = true
startretries     = 10
startsecs        = 2
stopsignal       = TERM
stopwaitsecs     = 20
killasgroup      = true
stopasgroup      = true
stdout_logfile   = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile   = /dev/stderr
stderr_logfile_maxbytes = 0
priority         = 10

[program:ct-core-daemon]
command          = /usr/local/bin/ct-server-core daemon
autostart        = true
autorestart      = true
startretries     = 10
startsecs        = 2
stopsignal       = TERM
stopwaitsecs     = 20
killasgroup      = true
stopasgroup      = true
stdout_logfile   = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile   = /dev/stderr
stderr_logfile_maxbytes = 0
priority         = 50
environment      = RUST_LOG="info,ct_server_core=info"
