# docker/panel/supervisord.conf
[supervisord]
nodaemon = true
logfile  = /dev/stdout
logfile_maxbytes = 0
pidfile  = /tmp/supervisord.pid
user     = root

# Post-runtime-swap layout: nginx + php-fpm decommissioned.
# FrankenPHP runs as a single in-process program (Caddy + PHP
# worker pool). The other long-lived processes (queue worker,
# scheduler, ct-server-core daemon) still need supervisord; the
# in-process Caddy doesn't replace them.
#
# Octane manages the FrankenPHP worker pool itself — this
# `frankenphp` program just keeps the parent Octane process
# resident; supervisord restarts it if it crashes.

[program:frankenphp]
# `octane:start` reads config/octane.php for server / host / port /
# workers defaults; environment overrides via OCTANE_* let
# operators tune without re-rolling the image. The --no-interaction
# flag suppresses the "are you sure?" prompts in CI / boot paths.
command          = /usr/local/bin/php /var/www/html/artisan octane:start --server=frankenphp --host=127.0.0.1 --port=9000 --workers=auto --no-interaction
autostart        = true
autorestart      = true
# A transient first-boot timing window where /run/cool-tunnel isn't
# yet writable, or where Octane's supervisor is still binding,
# would otherwise put this in FATAL after the supervisord default
# of 3 retries. 10 retries with a 5s bring-up window is generous.
startretries    = 10
startsecs       = 5
stdout_logfile   = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile   = /dev/stderr
stderr_logfile_maxbytes = 0
priority         = 10

[program:queue]
command          = php /var/www/html/artisan queue:work --sleep=1 --tries=3 --max-time=3600
autostart        = true
autorestart      = true
# A transient Redis blip during queue:work boot would otherwise put
# this program in FATAL — silently breaking the panel-side reload
# backstop. (v0.0.15.)
startretries    = 10
startsecs       = 5
stdout_logfile   = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile   = /dev/stderr
stderr_logfile_maxbytes = 0
priority         = 30

[program:scheduler]
# Laravel scheduler — kicks off traffic:rollup, quota:enforce, etc.
command          = /bin/sh -c 'while true; do php /var/www/html/artisan schedule:run --no-interaction; sleep 60; done'
autostart        = true
autorestart      = true
startretries    = 10
startsecs       = 5
stdout_logfile   = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile   = /dev/stderr
stderr_logfile_maxbytes = 0
priority         = 40

[program:ct-core-daemon]
# ct-server-core daemon — long-lived process that subscribes to
# Redis cool_tunnel:revocations and exposes a unix-socket JSON API
# the panel can call into. Restart-on-crash with a small backoff so
# transient Redis blips don't busy-loop. Unaffected by the
# FrankenPHP swap; the Rust side talks to the panel via the unix
# socket at /run/cool-tunnel/core.sock and never goes through the
# web tier.
command          = /usr/local/bin/ct-server-core daemon
autostart        = true
autorestart      = true
startretries     = 10
startsecs        = 2
stdout_logfile   = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile   = /dev/stderr
stderr_logfile_maxbytes = 0
priority         = 50
environment      = RUST_LOG="info,ct_server_core=info"
