# docker/panel/supervisord.conf
[supervisord]
nodaemon = true
logfile  = /dev/stdout
logfile_maxbytes = 0
pidfile  = /tmp/supervisord.pid
user     = root

[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0700

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Post-runtime-swap layout: nginx + php-fpm decommissioned.
# FrankenPHP runs as a single in-process program (Caddy + PHP
# worker pool). The other long-lived processes (queue worker,
# scheduler, ct-server-core daemon) still need supervisord; the
# in-process Caddy doesn't replace them.
#
# We invoke `/usr/local/bin/frankenphp run --config /etc/caddy/
# Caddyfile` directly — NOT through `php artisan octane:start
# --server=frankenphp`. Octane's CLI wrapper substitutes its own
# bundled stub Caddyfile (vendor/laravel/octane/src/Commands/
# stubs/Caddyfile), parameterised via CADDY_SERVER_* env vars,
# and silently discards any file we mounted at /etc/caddy/
# Caddyfile. Running FrankenPHP directly keeps our Caddyfile's
# anti-fingerprint customisations intact AND still loads
# Octane's per-request runtime hooks via the `worker { file ... }`
# block in the Caddyfile — that file is a thin shim that
# requires vendor/laravel/octane/bin/frankenphp-worker.php, so
# Filament + Livewire 3 still get the per-request container
# reset.
#
# Signal-forwarding discipline (added in the round-6 ops audit):
# every program has `stopsignal=TERM`, `stopwaitsecs=20`,
# `killasgroup=true`, `stopasgroup=true`. Pre-fix, supervisord
# defaulted to forwarding TERM only to the IMMEDIATE child;
# FrankenPHP's PHP worker subprocesses + Octane's tap-loop
# children would not see TERM, get SIGKILL'd via the cgroup
# when compose's 10 s grace window expired, and any in-flight
# request was lost. With killasgroup/stopasgroup the whole
# process group sees the signal at once; with stopwaitsecs=20
# Octane has time to drain workers cleanly before
# supervisord escalates.

[program:frankenphp]
# Config path is /etc/caddy-panel/ (NOT /etc/caddy/) — see the
# Dockerfile comment block above the matching COPY for the full
# rationale. tldr: docker-compose.yml mounts caddy_etc volume at
# /etc/caddy/ so the panel's CaddyfileGenerator can write the
# ACME Caddyfile for the dedicated ct-caddy container; that mount
# SHADOWS anything baked at /etc/caddy/ in the image, so the
# panel-internal FrankenPHP would load the ACME config instead of
# its own and never bind :9000. (post-v0.0.58 hotfix.)
command          = /usr/local/bin/frankenphp run --config /etc/caddy-panel/Caddyfile
autostart        = true
autorestart      = true
# A transient first-boot timing window where /run/cool-tunnel
# isn't yet writable, or where the worker pool is still binding,
# would otherwise put this in FATAL after the supervisord default
# of 3 retries. 10 retries with a 5s bring-up window is generous.
startretries    = 10
startsecs       = 5
stopsignal       = TERM
stopwaitsecs     = 20
killasgroup      = true
stopasgroup      = true
stdout_logfile   = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile   = /dev/stderr
stderr_logfile_maxbytes = 0
# MAX_REQUESTS bounds per-worker memory drift from any vendor
# package that doesn't reset cleanly between requests under
# Octane's long-lived workers (Filament forms, Livewire
# components have caused this in the wild). 500 mirrors the
# prior PHP-FPM `pm.max_requests` ceiling. Octane's
# vendor/laravel/octane/bin/frankenphp-worker.php reads this
# env var; default when unset is 1000. Operators on bigger
# boxes raise via docker-compose env override.
environment      = MAX_REQUESTS=500
priority         = 10

[program:queue]
command          = php /var/www/html/artisan queue:work --sleep=1 --tries=3 --max-time=3600
autostart        = true
autorestart      = true
# A transient Redis blip during queue:work boot would otherwise put
# this program in FATAL — silently breaking the panel-side reload
# backstop. (v0.0.15.)
startretries    = 10
startsecs       = 5
stopsignal       = TERM
stopwaitsecs     = 20
killasgroup      = true
stopasgroup      = true
stdout_logfile   = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile   = /dev/stderr
stderr_logfile_maxbytes = 0
priority         = 30

[program:messenger]
# Symfony Messenger worker — consumes the `cool_tunnel:messenger`
# Redis Streams transport and dispatches ReloadSingBox /
# ReloadServerConfig handlers. Runs side-by-side with
# [program:queue]: Laravel's queue worker still handles Filament
# internal jobs (mailables, etc.); Messenger handles the two
# reload-class messages explicitly routed there from
# MessengerServiceProvider's SendersLocator.
#
# --time-limit=3600 + --memory-limit=128M mirror the [program:queue]
# posture so both workers cycle on the same operator-visible cadence.
# supervisord auto-restarts the worker after exit (autorestart=true).
#
# startretries=10 + startsecs=5 same shape as [program:queue]:
# a Redis blip during boot must not put this program in FATAL.
# (v0.0.93, Phase 2 of the Symfony-infusion arc.)
command          = php /var/www/html/artisan messenger:consume async --time-limit=3600 --memory-limit=128M
autostart        = true
autorestart      = true
startretries    = 10
startsecs       = 5
stopsignal       = TERM
stopwaitsecs     = 20
killasgroup      = true
stopasgroup      = true
stdout_logfile   = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile   = /dev/stderr
stderr_logfile_maxbytes = 0
priority         = 31

[program:scheduler]
# Laravel scheduler — runs the active artisan maintenance jobs.
command          = /bin/sh -c 'while true; do php /var/www/html/artisan schedule:run --no-interaction; sleep 60; done'
autostart        = true
autorestart      = true
startretries    = 10
startsecs       = 5
stopsignal       = TERM
stopwaitsecs     = 20
killasgroup      = true
stopasgroup      = true
stdout_logfile   = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile   = /dev/stderr
stderr_logfile_maxbytes = 0
priority         = 40

[program:ct-core-daemon]
# ct-server-core daemon — long-lived process that exposes a
# unix-socket JSON API the panel can call into. Restart-on-crash
# with a small backoff so transient failures don't busy-loop.
# Unaffected by the FrankenPHP swap; the Rust side talks to the
# panel via the unix socket at /run/cool-tunnel/core.sock and
# never goes through the web tier.
command          = /usr/local/bin/ct-server-core daemon
autostart        = true
autorestart      = true
startretries     = 10
startsecs        = 2
stopsignal       = TERM
stopwaitsecs     = 20
killasgroup      = true
stopasgroup      = true
stdout_logfile   = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile   = /dev/stderr
stderr_logfile_maxbytes = 0
priority         = 50
environment      = RUST_LOG="info,ct_server_core=info"
