# docker/panel/supervisord.conf
[supervisord]
nodaemon = true
logfile  = /dev/stdout
logfile_maxbytes = 0
pidfile  = /tmp/supervisord.pid
user     = root

[program:php-fpm]
command          = php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.d/zz-pool.conf
autostart        = true
autorestart      = true
# startretries / startsecs override the supervisord defaults of (3 / 1)
# which would put the program in FATAL after three rapid crashes —
# e.g. a transient first-boot timing window where /run/cool-tunnel
# isn't yet writable. 10 retries with a 5s bring-up window is
# generous enough to ride out a slow VPS first-boot without ever
# losing the process for the lifetime of the container. (v0.0.15.)
startretries    = 10
startsecs       = 5
stdout_logfile   = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile   = /dev/stderr
stderr_logfile_maxbytes = 0
priority         = 10

[program:nginx]
command          = nginx -g "daemon off;"
autostart        = true
autorestart      = true
startretries    = 10
startsecs       = 5
stdout_logfile   = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile   = /dev/stderr
stderr_logfile_maxbytes = 0
priority         = 20

[program:queue]
command          = php /var/www/html/artisan queue:work --sleep=1 --tries=3 --max-time=3600
autostart        = true
autorestart      = true
# A transient Redis blip during queue:work boot would otherwise put
# this program in FATAL — silently breaking the panel-side reload
# backstop. (v0.0.15.)
startretries    = 10
startsecs       = 5
stdout_logfile   = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile   = /dev/stderr
stderr_logfile_maxbytes = 0
priority         = 30

[program:scheduler]
# Laravel scheduler — kicks off traffic:rollup, quota:enforce, etc.
command          = /bin/sh -c 'while true; do php /var/www/html/artisan schedule:run --no-interaction; sleep 60; done'
autostart        = true
autorestart      = true
startretries    = 10
startsecs       = 5
stdout_logfile   = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile   = /dev/stderr
stderr_logfile_maxbytes = 0
priority         = 40

[program:ct-core-daemon]
# ct-server-core daemon — long-lived process that subscribes to
# Redis cool_tunnel:revocations and exposes a unix-socket JSON API
# the panel can call into. Restart-on-crash with a small backoff so
# transient Redis blips don't busy-loop.
command          = /usr/local/bin/ct-server-core daemon
autostart        = true
autorestart      = true
startretries     = 10
startsecs        = 2
stdout_logfile   = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile   = /dev/stderr
stderr_logfile_maxbytes = 0
priority         = 50
environment      = RUST_LOG="info,ct_server_core=info"
